# releases/mh_tunnel/v6.0.0/rc.local
#!/bin/bash
##!/bin/sh -e
cd /root/mh_tunnel
nohup /root/mh_tunnel/mh_tunnel &
exit 0
